# rootfs/service/vusaged/run
#!/bin/sh
exec /var/vpopmail/bin/vusaged /var/vpopmail/etc/vusaged.conf 2>&1
