# rootfs/service/dovecot/run
#!/usr/bin/sh
exec /usr/sbin/dovecot -F 2>&1
