# rootfs/service/fcron/run
#!/usr/bin/sh
exec /usr/sbin/fcron -f -y 2>&1
